# config/mautrix-whatsapp.service
[Unit]
Description=Mautrix-WhatsApp bridge
After=network.target matrix-synapse.service postgresql.service

[Service]
Type=simple
Restart=always
RestartSec=5
User=matrix
Group=nogroup
WorkingDirectory=/home/matrix-ai
ExecStart=/home/matrix-ai/bin/mautrix-whatsapp -c /home/matrix-ai/config/mautrix-whatsapp/config.yaml --ignore-foreign-tables
Environment=HOME=/home/matrix-ai
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/home/matrix-ai
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target